# deploy/ops/vehicle-recycle-system-logrotate.service
[Unit]
Description=Rotate vehicle-recycle-system logs

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.d/vehicle-recycle-system
